# build/absinthe/cli/data/common/corona2/com.chronic-dev.greenpois0n.rocky-racoon.list
/.
/System
/System/Library
/System/Library/LaunchDaemons
/System/Library/LaunchDaemons/jb.plist
/tmp
/tmp/rocky_racoon_install
/tmp/rocky_racoon_install/bspatch
/tmp/rocky_racoon_install/payloads.tar.bz2
/tmp/rocky_racoon_install/racoon.bdiff
/usr
/usr/libexec
/usr/libexec/dirhelper
/var
/var/audit
/var/audit/rocky-racoon
/var/audit/rocky-racoon/crazeles
/var/audit/rocky-racoon/remounter
